# Systemd/empty-trash/empty-trash.service
[Unit]
Description=Empty trashed files older than 7 days for all users
ConditionPathExists=/usr/bin/trash-empty

[Service]
Type=oneshot
ExecStart=/usr/bin/trash-empty -f --all-users 7
User=root
